# backup/crontab
# GameForge Production Backup Crontab
# Automated backup schedule for production environment
# Runs daily at 2:00 AM UTC with error notifications

# Daily backup at 2:00 AM UTC
0 2 * * * /backup/scripts/backup.sh >> /backup/logs/cron.log 2>&1

# Weekly cleanup and maintenance on Sundays at 3:00 AM UTC
0 3 * * 0 /backup/scripts/maintenance.sh >> /backup/logs/cron.log 2>&1

# Health check every hour
0 * * * * echo "$(date): Backup service healthy" >> /backup/logs/health.log
